# crontab/factory-monitor
* * * * * root /var/lib/factory-monitor/monitor.py & >/dev/null 2>&1
* * * * * root (sleep 30; /var/lib/factory-monitor/monitor.py & >/dev/null 2>&1)

